# user/.config/systemd/user/checkupdates.timer
[Unit]
Description=Run checkupdate twice perday

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target